# bin/launchd/com.secondbrain.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    Second Brain — nightly sources→wiki sync.
    Runs brain-sync.sh every day at 02:00. launchd re-fires a missed
    StartCalendarInterval job when the Mac next wakes, so a sleeping
    machine still catches up. Managed via bin/brain-schedule.sh.
  -->
  <key>Label</key>
  <string>com.secondbrain.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/User/Tomi/second-brain/bin/brain-sync.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>2</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/User/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/User/Tomi/second-brain</string>

  <key>StandardOutPath</key>
  <string>/Users/User/Tomi/second-brain/.brain/launchd.sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/User/Tomi/second-brain/.brain/launchd.sync.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
